# docker/single/supervisord.conf
[program:ollama]
command=/usr/local/bin/ollama serve
environment=HOME="/root",OLLAMA_HOST="0.0.0.0:11434"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/ollama.log
stderr_logfile=/var/log/companion/ollama.err.log
startsecs=3

[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/15/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/postgres.log
stderr_logfile=/var/log/companion/postgres.err.log
startsecs=3

[program:qdrant]
command=/usr/local/bin/qdrant --config-path /opt/companion/config/qdrant.yaml
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/qdrant.log
stderr_logfile=/var/log/companion/qdrant.err.log
startsecs=2

[program:redis]
command=/usr/bin/redis-server /opt/companion/config/redis.conf
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/redis.log
stderr_logfile=/var/log/companion/redis.err.log
startsecs=1

[program:minio]
command=/usr/local/bin/minio server /var/lib/minio/data --console-address ":9001"
environment=MINIO_ROOT_USER="%(ENV_MINIO_ACCESS_KEY)s",MINIO_ROOT_PASSWORD="%(ENV_MINIO_SECRET_KEY)s"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/minio.log
stderr_logfile=/var/log/companion/minio.err.log
startsecs=2

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/nginx.log
stderr_logfile=/var/log/companion/nginx.err.log
startsecs=2

[program:hivemind]
command=/bin/bash -c "until pg_isready -h /var/run/postgresql -U postgres > /dev/null 2>&1; do sleep 1; done; /usr/local/bin/hivemind"
environment=DB_HOST="localhost",DB_PORT="5432",DB_NAME="postgres",DB_USER="supabase_admin",DB_PASSWORD="postgres",DB_SCHEMA="hivemind",JWT_SECRET="952572608172ccb759c0107b1ea873c1c7b3e0561753effc1b9a0be15a22c512",JWT_TTL_SECONDS="2592000",ENCRYPTION_SECRET="7fc09e1cee970e2ab7fa4880a525998277712c393de94aa42cc725d3bad8081b",VEXA_API_URL="http://localhost:8056",VEXA_ADMIN_API_URL="http://localhost:8057",VEXA_ADMIN_TOKEN="token",EMBEDDING_API_URL="http://localhost:11434",EMBEDDING_MODEL="nomic-embed-text",QDRANT_URL="http://localhost:6334",QDRANT_API_KEY="qdrant-dev-key",HOST="0.0.0.0",PORT="9100"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/companion/hivemind.log
stderr_logfile=/var/log/companion/hivemind.err.log
startsecs=3

[program:vexa-api-gateway]
command=/opt/companion/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8056
directory=/opt/companion/services/api-gateway
environment=ADMIN_API_URL="http://localhost:8057",BOT_MANAGER_URL="http://localhost:8080",TRANSCRIPTION_COLLECTOR_URL="http://localhost:8124",MCP_URL="http://localhost:18888",REDIS_URL="redis://localhost:6379/0",PUBLIC_BASE_URL="http://localhost:8056",DB_HOST="localhost",DB_PORT="5432",DB_NAME="postgres",DB_USER="supabase_admin",DB_PASSWORD="postgres",DB_SCHEMA="vexa"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/api-gateway.log
stderr_logfile=/var/log/vexa/api-gateway.err.log
startsecs=5

[program:vexa-admin-api]
command=/opt/companion/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8057
directory=/opt/companion/services/admin-api
environment=DB_HOST="localhost",DB_PORT="5432",DB_NAME="postgres",DB_USER="supabase_admin",DB_PASSWORD="postgres",DB_SCHEMA="vexa",ADMIN_API_TOKEN="token"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/admin-api.log
stderr_logfile=/var/log/vexa/admin-api.err.log
startsecs=3

[program:vexa-bot-manager]
command=/opt/companion/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/opt/companion/services/bot-manager
environment=REDIS_URL="redis://localhost:6379/0",TTS_SERVICE_URL="http://localhost:8002",DB_HOST="localhost",DB_PORT="5432",DB_NAME="postgres",DB_USER="supabase_admin",DB_PASSWORD="postgres",DB_SCHEMA="vexa",ADMIN_TOKEN="token",DEVICE_TYPE="remote",WHISPER_LIVE_URL="ws://localhost:9090/ws",STORAGE_BACKEND="minio",MINIO_ENDPOINT="localhost:9000",MINIO_ACCESS_KEY="vexa-access-key",MINIO_SECRET_KEY="vexa-secret-key",MINIO_BUCKET="vexa-recordings"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/bot-manager.log
stderr_logfile=/var/log/vexa/bot-manager.err.log
startsecs=5

[program:vexa-transcription-collector]
command=/opt/companion/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8124
directory=/opt/companion/services/transcription-collector
environment=DB_HOST="localhost",DB_PORT="5432",DB_NAME="postgres",DB_USER="supabase_admin",DB_PASSWORD="postgres",DB_SCHEMA="vexa",REDIS_HOST="localhost",REDIS_PORT="6379",ADMIN_TOKEN="token",STORAGE_BACKEND="minio",MINIO_ENDPOINT="localhost:9000",MINIO_ACCESS_KEY="vexa-access-key",MINIO_SECRET_KEY="vexa-secret-key",MINIO_BUCKET="vexa-recordings"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/transcription-collector.log
stderr_logfile=/var/log/vexa/transcription-collector.err.log
startsecs=5

[program:vexa-mcp]
command=/opt/companion/venv/bin/uvicorn main:app --host 0.0.0.0 --port 18888
directory=/opt/companion/services/mcp
environment=API_GATEWAY_URL="http://localhost:8056"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/mcp.log
stderr_logfile=/var/log/vexa/mcp.err.log
startsecs=3

[program:vexa-tts-service]
command=/opt/companion/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8002
directory=/opt/companion/services/tts-service
environment=OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",OPENAI_BASE_URL="https://api.openai.com"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vexa/tts-service.log
stderr_logfile=/var/log/vexa/tts-service.err.log
startsecs=3
